# automation/launchd/com.gamescaler.ai-coding-daily-fallback.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gamescaler.ai-coding-daily-fallback</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/mvbj0638/.ai_coding_survey_delivery_runtime/scripts/run_daily_fallback_and_delivery.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>11</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>11</integer>
      <key>Minute</key>
      <integer>45</integer>
    </dict>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Library/Frameworks/Python.framework/Versions/3.14/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/mvbj0638/.ai_coding_survey_delivery_runtime/data/feishu/launchd_daily_fallback.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mvbj0638/.ai_coding_survey_delivery_runtime/data/feishu/launchd_daily_fallback.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/mvbj0638/.ai_coding_survey_delivery_runtime</string>
</dict>
</plist>
